# etc/supervisor.d/gethotelcat.ini
[program:gethotelcat]
directory = /home/vnp/sources/hotel-crawler/crawler
command = /home/vnp/sources/hotel-crawler/env/bin/python /home/vnp/sources/hotel-crawler/crawler/crawler/run.py --spider=AgodaCat,BookingCat,HotelsCat,IvivuCat,LuxstayCat,MytourCat,ReddoorzCat,TravelokaCat,TripCat,TripiCat,TripadvisorCat,VnTripCat
autostart=true
autorestart=false
numprocs=1
priority = 900
redirect_stderr=true
stderr_logfile=var/log/%(program_name)s.log
stdout_logfile=var/log/%(program_name)s.log