# etc/supervisor.d/gethoteldetail.ini
[program:gethoteldetail]
directory = /home/vnp/sources/hotel-crawler/crawler
command = /home/vnp/sources/hotel-crawler/env/bin/python /home/vnp/sources/hotel-crawler/crawler/crawler/run.py --spider=AgodaDetail,BookingDetail,HotelsDetail,IvivuDetail,LuxstayDetail,MytourDetail,ReddoorzDetail,TravelokaDetail,TripDetail,TripiDetail,TripadvisorDetail,VnTripDetail
autostart=true
autorestart=false
numprocs=1
priority = 900
redirect_stderr=true
stderr_logfile=var/log/%(program_name)s.log
stdout_logfile=var/log/%(program_name)s.log