# etc/supervisor.d/tintuc_transport_nh.ini
[program:tintuc_transport_nh]
directory = /home/vnp/sources/hotel-crawler/crawler
command = /home/vnp/sources/hotel-crawler/env/bin/python /home/vnp/sources/hotel-crawler/crawler/crawler/run.py --transport=CookyDetail,ZingDetail    
autostart=true
autorestart=true
numprocs=1
priority = 900
redirect_stderr=true
stderr_logfile=var/log/%(program_name)s.log
stdout_logfile=var/log/%(program_name)s.log